# files/wikijs.service.j2
[Unit]
Description=Wiki.js
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node server
Restart=always
# Consider creating a dedicated user for Wiki.js here:
User={{ wikijs_unix_user }}
Environment=NODE_ENV=production
WorkingDirectory={{ wikijs_dir | expanduser }}

[Install]
WantedBy=multi-user.target